# data/gsignond.service
[Unit]
Description=Single SignOn Daemon
Requires=dbus.service

[Service]
ExecStart=/usr/bin/gsignond

[Install]
WantedBy=default.target
